# dht-cntl.service
[Unit]
Description=Transmission DHT Control

[Service]
Type=oneshot
ExecStart=/<path>/<to>/dht-cntl.sh

[Install]
WantedBy=multi-user.target
